# config/copyparty.service
[Unit]
Description=copyparty file server
After=network.target

[Service]
Type=simple
User=copyparty
Group=copyparty
WorkingDirectory=/opt/copyparty
EnvironmentFile=/opt/copyparty/config/secrets.env
ExecStart=/usr/bin/python3 /opt/copyparty/copyparty-sfx.py -c /opt/copyparty/config/copyparty.conf --xau f,j,/opt/copyparty/hooks/xau-hook.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
